# homelab/systemd/ansible-node-state.service
[Unit]
Description=Ansible node-state playbook
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ansible
ExecStart=/usr/local/bin/run-node-state.sh
StandardOutput=journal
StandardError=journal
